# 85-vaiopower.rules.in
# udev rules for vaiopower

SUBSYSTEM=="rfkill",   ACTION=="change", ATTR{type}=="wlan", ATTR{state}=="[01]", RUN+="###SBINDIR###/vaiopower rfkill"
ENV{EVENT}=="dock", KERNEL=="dock.0", SUBSYSTEM=="platform", RUN+="###SBINDIR###/vaiopower dock"

