# init-scripts/upstart/vaiopower.conf
# vaiopower - power management tool for Sony Vaio
#

description		"vaiopower: power management tool for Sony Vaio"

start on filesystem

task
exec /usr/sbin/vaiopower boot
